# navi-dashboard.service
[Unit]
Description=Navi Dashboard
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/.openclaw/workspace/navi-dashboard
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=/home/ubuntu/.openclaw/workspace/navi-dashboard/.env.local

[Install]
WantedBy=multi-user.target
